# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/dreamer.conf
[program:dreamer]
command=/opt/ai-dock/bin/supervisor-dreamer.sh
process_name=="%(program_name)s_%(process_num)s"
numprocs=%(ENV_GPU_COUNT)s
directory=/opt/AI-Horde-Worker
priority=2000
autostart=true
startsecs=0
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/dreamer-%(process_num)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/dreamer-%(process_num)s.error.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=1
environment=PROC_NAME="%(program_name)s",PROC_NUM=%(process_num)d
